# util/node-scimitar.conf
#!upstart

# Scimitar upstart script for RedHat EL 6 (CentOS 6)

description "Scimitar - node.js reconnoiter console"
author      "michal@taborsky.cz"

start on stopped rc RUNLEVEL=[2345]
stop on starting runlevel [016]

respawn                # restart when job dies
respawn limit 5 60     # give up restart after 5 respawns in 60 seconds

script
  echo $$ > /var/run/scimitar.pid
  export HOME="/root"  
  exec node /opt/scimitar/scimitar.js >> /var/log/node-scimitar.log 2>&1
end script

pre-stop script
    rm /var/run/scimitar.pid
end script
